# s6/etc/services.d/nullmailer/finish
#!/usr/bin/with-contenv bash

s6-echo "Stopping nullmailer"
killall -9 nullmailer-send
